# robot2.service
[Unit]
Description=Robot 2 vision

[Service]
Type=simple
ExecStart=/bin/bash /home/zegion/autorun.sh

[Install]
WantedBy=multi-user.target

